# infra/systemd/jbnu-sugang-helper-dr-backup.timer
[Unit]
Description=Daily JBNU Sugang Helper encrypted DR backup schedule

[Timer]
OnCalendar=*-*-* 03:17:00
Persistent=true
RandomizedDelaySec=10m
Unit=jbnu-sugang-helper-dr-backup.service

[Install]
WantedBy=timers.target
